# scripts/hello-flask-gunicorn.conf
[program:hello-world-gunicorn]
command = /bin/hello-flask-gunicorn
autostart=true
autorestart=true
stderr_logfile=/var/log/start_gunicorn.err.log
stdout_logfile=/var/log/start_gunicorn.out.log
